# ryndns.service
[Unit]
Description=Update DNS records

[Service]
Type=oneshot
ExecStart=/opt/ryndns/ryndns
